# edge_api_server/deploy/rist-edge-api.service
[Unit]
Description=RIST Edge API Server (FastAPI/Uvicorn) on 0.0.0.0:8000
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rist
Group=rist

WorkingDirectory=/home/rist/ritas/edge_api_server
Environment=RIST_ENV=production
# DB 접속 정보 등 환경/비밀값은 git 추적 밖의 파일에서 읽는다.
# 앞의 '-' 는 파일이 없어도 기동 실패하지 않게 한다(설치 시 생성).
EnvironmentFile=-/home/rist/ritas/edge.env
# 로깅: 레벨/포맷/파일 경로. 파일 회전은 RIST_LOG_MAX_BYTES, RIST_LOG_BACKUP_COUNT 로 조정.
Environment=RIST_LOG_LEVEL=INFO
Environment=RIST_LOG_FORMAT=text
Environment=RIST_LOG_DIR=/home/rist/ritas/edge_api_server/data/logs
ExecStart=/home/rist/ritas/.venv/bin/python -m app.run

Restart=always
RestartSec=5
# 종료 시 정상 종료를 위해 SIGINT 사용(Uvicorn graceful shutdown)
KillSignal=SIGINT

# 보안 강화
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
