# edge_api_server/deploy/rist-edge-worker.service
[Unit]
Description=RIST Edge Report Worker (LLM report generator)
After=network-online.target docker.service rist-edge-api.service
# vLLM 은 docker-compose 컨테이너(vllm-gemma4-e4b, http://127.0.0.1:8001)로 구동된다.
# 컨테이너가 준비되기 전이라도 worker 는 기동하며, LLM 이 없으면 해당 작업만 실패한다.
# 따라서 docker 를 선택 의존(Wants)으로 두어 순서만 맞춘다.
Wants=network-online.target docker.service

[Service]
Type=simple
User=rist
Group=rist

WorkingDirectory=/home/rist/ritas/edge_api_server
Environment=RIST_ENV=production
# API 서비스와 동일한 환경/비밀값 파일을 공유한다(git 추적 밖).
EnvironmentFile=-/home/rist/ritas/edge.env
# 로깅: API 서비스와 동일 정책. worker 로그는 별도 파일로 분리한다.
Environment=RIST_LOG_LEVEL=INFO
Environment=RIST_LOG_FORMAT=text
Environment=RIST_LOG_FILE=/home/rist/ritas/edge_api_server/data/logs/worker.log
ExecStart=/home/rist/ritas/.venv/bin/python -m app.report_worker

Restart=always
RestartSec=5

# 보안 강화
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
