# edge_api_server/deploy/rist-edge-worker.service
[Unit]
Description=RIST Edge Report Worker (LLM report generator)
After=network-online.target docker.service rist-vllm.service rist-edge-api.service
# vLLM 은 rist-vllm.service(docker compose, vllm-gemma4-e4b, http://127.0.0.1:8001)로 구동된다.
# worker 를 start 하면 Wants 로 vLLM 스택도 함께 기동된다(After 로 순서 보장).
# 단, 선택 의존(Wants)이라 vLLM 이 떠 있지 않거나 준비 전이라도 worker 는 기동하며,
# LLM 이 없으면 해당 작업의 보조 문안만 규칙 기본값으로 대체되고 보고서는 완성된다.
Wants=network-online.target rist-vllm.service

[Service]
Type=simple
User=rist
Group=rist

WorkingDirectory=/home/rist/ritas/edge_api_server
Environment=RIST_ENV=production
# API 서비스와 동일한 환경/비밀값 파일을 공유한다(git 추적 밖).
EnvironmentFile=-/home/rist/ritas/edge.env
# 로깅: API 서비스와 동일 정책. worker 로그는 별도 파일로 분리한다.
Environment=RIST_LOG_LEVEL=INFO
Environment=RIST_LOG_FORMAT=text
Environment=RIST_LOG_FILE=/home/rist/ritas/edge_api_server/data/logs/worker.log
ExecStart=/home/rist/ritas/.venv/bin/python -m app.report_worker

Restart=always
RestartSec=5

# 보안 강화
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
